# private_dot_config/systemd/user/sway-session.target
[Unit]
Description=Sway Graphical Session
